# web3-auth-svc.service
[Unit]
Description=BlockHost Web3 Auth Service (Cardano)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/web3-auth-svc
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
